# scripts/ssh-tunnel-jetson-robot.service
[Unit]
Description=SSH Reverse Tunnel to VPS for Jetson Robot Project
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=jetson
Restart=always
RestartSec=10
StartLimitInterval=0

# Reverse SSH tunnel: VPS:2223 -> localhost:22 (SSH access)
#                      VPS:8554 -> localhost:8554 (RTSP camera stream)
# This allows:
#   - Accessing Jetson via: ssh -p 2223 jetson@81.200.157.230
#   - Viewing camera via: rtsp://81.200.157.230:8554/camera
ExecStart=/usr/bin/ssh -N \
    -o ServerAliveInterval=30 \
    -o ServerAliveCountMax=3 \
    -o ExitOnForwardFailure=yes \
    -o StrictHostKeyChecking=accept-new \
    -i /home/jetson/.ssh/id_ed25519_vps \
    -R 2223:localhost:22 \
    -R 8554:localhost:8554 \
    root@81.200.157.230

[Install]
WantedBy=multi-user.target
